# builders/packing/debian/debian/crtmpserver.postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ] ; then
	update-rc.d crtmpserver remove >/dev/null
	rm -rf /var/lib/crtmpserver
	rm -rf /var/log/crtmpserver
	deluser crtmpserver
fi

#DEBHELPER#
