# builders/packing/debian/debian/crtmpserver.postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ] ; then
	update-rc.d crtmpserver remove >/dev/null
	rmdir /var/lib/crtmpserver
fi
