# builders/packing/debian/debian/rtmpd.postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ] ; then
	update-rc.d rtmpd remove >/dev/null
	rm -rf /var/lib/rtmpd
	rm -rf /var/log/rtmpd
	deluser rtmpd
fi

#DEBHELPER#
