# builders/packing/debian/debian/rtmpd.prerm
#!/bin/sh

set -e

if [ -x "/etc/init.d/rtmpd" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d rtmpd stop || exit $?
	else
		/etc/init.d/rtmpd stop || exit $?
	fi
fi

