# builders/packing/debian/debian/rtmpd.postinst
#!/bin/sh

set -e

USER="rtmpd"
GROUP="$USER"

chown $USER:$GROUP /var/lib/rtmpd
chown $USER:$GROUP /var/log/rtmpd

if [ -x "/etc/init.d/rtmpd" ]; then
	update-rc.d rtmpd defaults >/dev/null
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d rtmpd start || exit $?
	else
		/etc/init.d/rtmpd start || exit $?
	fi
fi

#DEBHELPER# 
